# sftpgo.service.example
[Unit]
Description=SFTP golang docker server

[Service]
ExecStart=docker-compose -f ${APATH} up
