# sftpgo.service.example
[Unit]
Description=SFTP golang docker server

[Service]
EnvironmentFile=${AWD}/.env.private
ExecStart=docker-compose -f ${AWD}/docker-compose.yml up
